# dwl/scripts/wallpaper.cron
* * * * *  /home/thiago/bin/wallpaper.sh  